# contrib/init/cryptopesetad.conf
description "CryptoPeseta Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env CRYPTOPESETAD_BIN="/usr/bin/cryptopesetad"
env CRYPTOPESETAD_USER="cryptopeseta"
env CRYPTOPESETAD_GROUP="cryptopeseta"
env CRYPTOPESETAD_PIDDIR="/var/run/cryptopesetad"
# upstart can't handle variables constructed with other variables
env CRYPTOPESETAD_PIDFILE="/var/run/cryptopesetad/cryptopesetad.pid"
env CRYPTOPESETAD_CONFIGFILE="/etc/cryptopeseta/cryptopeseta.conf"
env CRYPTOPESETAD_DATADIR="/var/lib/cryptopesetad"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # cryptopesetad will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$CRYPTOPESETAD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run cryptopesetad."
        echo "The setting must appear in $CRYPTOPESETAD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"CryptoPeseta Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$CRYPTOPESETAD_PIDDIR"
    chmod 0755 "$CRYPTOPESETAD_PIDDIR"
    chown $CRYPTOPESETAD_USER:$CRYPTOPESETAD_GROUP "$CRYPTOPESETAD_PIDDIR"
    chown $CRYPTOPESETAD_USER:$CRYPTOPESETAD_GROUP "$CRYPTOPESETAD_CONFIGFILE"
    chmod 0660 "$CRYPTOPESETAD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$CRYPTOPESETAD_PIDFILE" \
    --chuid $CRYPTOPESETAD_USER:$CRYPTOPESETAD_GROUP \
    --exec "$CRYPTOPESETAD_BIN" \
    -- \
    -pid="$CRYPTOPESETAD_PIDFILE" \
    -conf="$CRYPTOPESETAD_CONFIGFILE" \
    -datadir="$CRYPTOPESETAD_DATADIR" \
    -disablewallet \
    -daemon

